# docker/systemd/floway-server.container.example
[Container]
ContainerName=floway-server
Environment="ADMIN_KEY=<admin-secret>" FLOWAY_DB_PATH=/data/floway.db FLOWAY_FILES_DIR=/data/files PORT=8788 RUNTIME_LOCATION=LOCAL
HealthCmd=["node", "-e", "fetch('http://127.0.0.1:8788/api/health').then(r=>process.exit(r.ok?0:1)).catch(()=>process.exit(1))"]
HealthInterval=10s
HealthOnFailure=stop
HealthRetries=12
HealthStartPeriod=15s
HealthTimeout=5s
Image=ghcr.io/menci/floway-server:latest
Mount=type=volume,source=floway-data.volume,destination=/data
Notify=healthy
Pod=floway.pod

# Auto-update
Pull=newer

[Service]
Restart=always
